# tunneler/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; just some examples, may figure out other stuff to put here

[program:wormhole]
command=twist3 transitrelay --port=tcp:%(ENV_WORMHOLE_PORT)s --blur-usage=100
autostart=%(ENV_WORMHOLE_AUTOSTART)s
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:hostyoself]
command=hostyoself relay
    --url "%(ENV_HOSTYOSELF_URL)s"
    --port %(ENV_HOSTYOSELF_PORT)s
    --google-site-verification-id "%(ENV_HOSTYOSELF_GOOGLE_SITE_VER_ID)s"
    --umami-url "%(ENV_HOSTYOSELF_UMAMI_URL)s"
    --umami-id  "%(ENV_HOSTYOSELF_UMAMI_ID)s"
autostart=%(ENV_HOSTYOSELF_AUTOSTART)s
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:croc]
command=/usr/local/bin/croc_relay_entrypoint.sh
autostart=%(ENV_CROC_AUTOSTART)s
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true

[program:wireproxy]
command=/usr/bin/wireproxy --config "%(ENV_WIREPROXY_CONF)s"
autostart=%(ENV_WIREPROXY_AUTOSTART)s
autorestart=true
startsecs=0
startretries=0
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
redirect_stderr=true
